# docker_service_watch.service
[Unit]
Description=Docker Service Watch Dog

[Service]
Type=simple
ExecStart=XXXX
User=root
WorkingDirectory=YYYY

[Install]
WantedBy=default.target